# basic-phpsite/config/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm7 -F

[program:nginx]
command=nginx -g 'daemon off;'
